# service_templates/org.dyyni.gs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.dyyni.gs</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/full/absolute/path/to/your/gs</string>
        <string>auto</string>
        <string>--timeout</string>
        <string>5m</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <!-- only run once, don't restart -->
    <key>LaunchOnlyOnce</key>
    <true/>
    
    <!-- allow potential hangs when testing network connection -->
    <key>TimeOut</key>
    <integer>30</integer>
</dict>
</plist>

